# syhub/syhub.service
[Unit]
Description=syhub

[Service]
ExecStart=

[Install]
WantedBy=multi-user.target